# cpu_limit.cron
*/5 * * * * cd /root/cpu/&&/bin/bash ./cpu_core 
